# systemd/eda.service
[Unit]
Description=EDA schematic/PCB web server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/epentland/ai/canopy/eda
ExecStart=/home/epentland/ai/canopy/eda/zig-out/bin/eda serve --project-dir /home/epentland/ai/canopy/eda/projects/designs --port 7050
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
