# systemd/netlisp.service
[Unit]
Description=Netlisp schematic/PCB web server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/epentland/ai/canopy/eda
# Pin the DEPLOYED binary to ReleaseSafe (see the rationale next to
# standardOptimizeOption in build.zig). The server parses untrusted input
# (.sexp designs, imported KiCad boards, HTTP/JSON params); in a safety-off
# build any remaining unguarded cast/overflow is silent UB — a wrong board —
# whereas ReleaseSafe turns it into a panic, and Restart=on-failure below makes
# that a ~2-second blip. A plain `zig build` ships Debug; this rebuilds
# ReleaseSafe before serving. Dev/test/mutation stay Debug (fast iteration),
# which is why the pin lives here at the deploy point, not in build.zig's
# default. Login shell so `zig` is on PATH (mirrors the nightly unit). The
# guardian `all` cache + incremental link keep an unchanged-tree restart fast.
ExecStartPre=/bin/sh -lc 'cd /home/epentland/ai/canopy/eda && zig build -Doptimize=ReleaseSafe'
ExecStart=/home/epentland/ai/canopy/eda/zig-out/bin/netlisp serve --project-dir /home/epentland/ai/canopy/eda/projects/designs --port 7050
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
